# systemd/llama-server.service
[Unit]
Description=llama.cpp server (gemma4-e2b vision)
After=network.target

[Service]
Type=simple
User=jsligar
ExecStart=/opt/llama.cpp/build/bin/llama-server \
  -m /opt/models/gemma4-e2b-q4km-compare/sowilow/gemma-4-e2b-it-q4_k_m.gguf \
  --mmproj /opt/models/gemma4-e2b-ud-q4xl-fresh/mmproj-F16.gguf \
  --host 0.0.0.0 \
  --port 8424 \
  -ngl 99 \
  -c 32768 \
  -t 6 -tb 6 \
  --threads-http 1 \
  -np 1 \
  --no-cont-batching \
  -b 1024 -ub 1024 \
  --flash-attn on \
  --cache-type-k q8_0 \
  --cache-type-v q8_0 \
  --prio 2 --prio-batch 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
